# scripts/bin/kroesclass/stagingroot/etc/systemd/system/cockpit.socket.d/listen.conf
[Socket]
ListenStream=
ListenStream=9999
